# airootfs/etc/skel/.config/autostart/ERICGNOME.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Exec=ericgnome.sh
Name=EricGNOME
Hidden=false
